# includes.container/etc/profile.d/add_libvirt_group.sh
#!/bin/sh
usermod -aG libvirt $USER
usermod -aG libvirt-qemu $USER
usermod -aG kvm $USER
